# meeting-followups/meeting-followups.timer
[Unit]
Description=Daily meeting-to-task follow-up reconciliation

[Timer]
# Fire daily at ~07:00 local, before the workday starts.
OnCalendar=*-*-* 07:00:00
# Spread the wake-up a little so it doesn't fight other 7am jobs.
RandomizedDelaySec=5min
# Catch up on a run missed while the machine was asleep/off.
Persistent=true

[Install]
WantedBy=timers.target
